# deploy/thg-frontend.service
# /etc/systemd/system/thg-frontend.service
# Install: sudo systemctl daemon-reload && sudo systemctl enable thg-frontend
# Requires Node.js on server: curl -fsSL https://deb.nodesource.com/setup_20.x | sudo -E bash - && sudo apt-get install -y nodejs

[Unit]
Description=THG Lead Intelligence — Next.js Frontend
After=network.target thg-scraper.service
Wants=network-online.target

[Service]
Type=simple

User=thg-scraper
Group=thg-scraper

# Standalone server — built by CI with `output: 'standalone'`
WorkingDirectory=/opt/thg-scraper/frontend-standalone
# /usr/bin/env resolves 'node' from PATH — works regardless of NodeSource install path
ExecStart=/usr/bin/env node server.js

Environment=API_URL=http://127.0.0.1:8080
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3000

Restart=on-failure
RestartSec=10
StartLimitInterval=120
StartLimitBurst=5

MemoryMax=512M
TasksMax=256

NoNewPrivileges=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=thg-frontend

[Install]
WantedBy=multi-user.target
