# deploy/thg-agent-brain.service
# /etc/systemd/system/thg-agent-brain.service
# Install: sudo systemctl daemon-reload && sudo systemctl enable thg-agent-brain

[Unit]
Description=THG Agent Brain (LLM Planner Sidecar)
After=network.target
Before=thg-scraper.service

[Service]
Type=simple

# ── Run as same user as scraper ───────────────────────────────────────────
User=thg-scraper
Group=thg-scraper

# ── Working directory and binary ──────────────────────────────────────────
WorkingDirectory=/opt/thg-agent-brain
ExecStart=/opt/thg-agent-brain/venv/bin/python brain.py

# ── Config ────────────────────────────────────────────────────────────────
# Bind to localhost only — scraper connects via loopback
Environment=AGENT_BRAIN_HOST=127.0.0.1
Environment=AGENT_BRAIN_PORT=8091
Environment=AGENT_BRAIN_PLANNER_MODEL=gpt-4o-mini
Environment=AGENT_BRAIN_PLANNER_TIMEOUT_S=5

# Reads OPENAI_API_KEY (and other secrets) from scraper env file
EnvironmentFile=/etc/thg-scraper/env

# ── Crash recovery ────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=5
StartLimitInterval=120
StartLimitBurst=10

# ── Resource limits (very lightweight service) ────────────────────────────
MemoryMax=256M
TasksMax=64

# ── Hardening (principle of least privilege) ──────────────────────────────
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/thg-agent-brain

# Prevent the process from gaining new capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# Network: allow only outbound (no raw sockets)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# No kernel module loading
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true

# ── Logging ───────────────────────────────────────────────────────────────
# Logs go to journald: sudo journalctl -u thg-agent-brain -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thg-agent-brain

[Install]
WantedBy=multi-user.target
