# deploy/thg-worker.service
# /etc/systemd/system/thg-worker.service
# Install: sudo systemctl daemon-reload && sudo systemctl enable thg-worker

[Unit]
Description=THG Lead Intelligence — Background Worker
After=network.target thg-scraper.service
Wants=network-online.target

[Service]
Type=simple

User=thg-scraper
Group=thg-scraper

WorkingDirectory=/opt/thg-scraper
ExecStart=/opt/thg-scraper/thg-worker

EnvironmentFile=/etc/thg-scraper/env

Restart=on-failure
RestartSec=15
StartLimitInterval=120
StartLimitBurst=5

MemoryMax=1G
TasksMax=256

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/thg-scraper/data /var/run/docker.sock /tmp
SupplementaryGroups=docker
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=thg-worker

[Install]
WantedBy=multi-user.target
