# deploy/thg-api.service
# /etc/systemd/system/thg-api.service
# Install: sudo systemctl daemon-reload && sudo systemctl enable thg-api

[Unit]
Description=THG Lead Intelligence — API Server
After=network.target
Wants=network-online.target

[Service]
Type=simple

User=thg-scraper
Group=thg-scraper

WorkingDirectory=/opt/thg-scraper
ExecStart=/opt/thg-scraper/thg-api

# Shared secrets (DB_PATH, OPENAI_API_KEY, etc.)
EnvironmentFile=/etc/thg-scraper/env
# API server runs on 8081 so it doesn't conflict with cmd/scraper on 8080
Environment=API_PORT=8081

Restart=on-failure
RestartSec=10
StartLimitInterval=120
StartLimitBurst=5

MemoryMax=512M
TasksMax=256

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/thg-scraper/data
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=thg-api

[Install]
WantedBy=multi-user.target
